# op5build/naemon.service.in
[Unit]
Description=Naemon Monitoring Daemon
Documentation=http://naemon.org/documentation
After=network.target

[Service]
Type=forking
PIDFile=@lockfile@
PermissionsStartOnly=true
ExecStartPre=/usr/bin/sh /etc/sysconfig/naemon
ExecStartPre=-/usr/bin/mkdir --parents @TMPDIR@
ExecStartPre=/usr/bin/chown --recursive monitor:apache @TMPDIR@
ExecStartPre=su monitor --login --shell=/bin/sh --command="@bindir@/naemon --verify-config --precache-objects @pkgconfdir@/naemon.cfg"
ExecStart=@bindir@/naemon --daemon @pkgconfdir@/naemon.cfg
ExecReload=/bin/kill -HUP $MAINPID
User=monitor
Group=apache
StandardOutput=journal
StandardError=inherit

[Install]
WantedBy=multi-user.target
Alias=monitor.service
